# systemd/teuthology-node-cleanup@.timer
[Unit]
Description=Run Teuthology Node Cleanup

[Timer]
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
